# config/supervisor/supervisord_docker.conf
[unix_http_server]
file=/tmp/supervisor.sock    ; the path to the socket file
chmod=777                    ; socket file mode
chown=root:root              ; socket file uid:gid owner

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
umask=022                    ; process file creation umask; default 022
identifier=supervisor        ; supervisord identifier, default is 'supervisor'
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:gunicorn]
command=gunicorn --paste /app/config/environment/docker.ini --reload
user=root
numprocs=1
stdout_logfile=/var/log/ckan/app.log
stderr_logfile=/var/log/ckan/app.log
autostart=true
autorestart=true
startsecs=1

#[program:ckan_gather_consumer]
#command=paster --plugin=ckanext-harvest harvester gather_consumer --config=/app/config/docker.ini
#user=root
#numprocs=1
#stdout_logfile=/var/log/ckan/gather_consumer.log
#stderr_logfile=/var/log/ckan/gather_consumer.log
#autostart=true
#autorestart=true
#startsecs=10
#
#[program:ckan_fetch_consumer]
#command=paster --plugin=ckanext-harvest harvester fetch_consumer --config=/app/config/docker.ini
#user=root
#numprocs=1
#stdout_logfile=/var/log/ckan/fetch_consumer.log
#stderr_logfile=/var/log/ckan/fetch_consumer.log
#autostart=true
#autorestart=true
#startsecs=10

[eventlistener:crashmail]
command=crashmail -s "/usr/sbin/sendmail -S mailhog:1025" -p ckan_gather_consumer -p ckan_fetch_consumer -m team-gerwin@list.dkd.de -o "Transparenzplatform (tpp.rlp.de)"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_UNKNOWN,PROCESS_STATE_FATAL,TICK_60
stderr_logfile=/var/log/ckan/supervisord_stderr.log
stdout_logfile=/var/log/ckan/supervisord_stdout.log